# github_backup/crontab
# Run the backup script daily at midnight
0 0 * * * /usr/local/bin/python /app/github_backup.py >> /var/log/cron.log 2>&1
