# radicle_push/crontab
# Run the Radicle push script daily at 1 AM
0 1 * * * python /app/radicle_push.py >> /var/log/cron.log 2>&1
